# crontab
* */1   * * *     root   ntpdate-debian
30 11   * * 1-5/  root   playmusic.sh
30 12   * * 1-5/  root   pkill mplayer  && playring.sh